# Defense Module.conf
# Dynamic A-D Ports for Applications (DAD Ports for App)
#
# addon for gufw - ufw

start on (starting ufw)

stop on runlevel [!023456]

exec "/opt/DAD Ports for App/Defense Module"
